# commons-server.service
[Unit]
Description=Commons Game Server
After=network.target

[Service]
Type=simple
WorkingDirectory=/mnt/data/commons-server
ExecStart=/home/clungus/.bun/bin/bun run src/index.ts
Restart=always
RestartSec=10
SyslogIdentifier=commons-server
Environment=NODE_ENV=production

[Install]
WantedBy=default.target
